# linux/aide/aidecheck.timer
[Unit]
Description=Daily AIDE check timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target